# systemd/user/sync-tigerc.timer
[Unit]
Description=Run sync TigerC

[Timer]
OnCalendar=*-*-* 00:00:2
WakeSystem=true

[Install]
WantedBy=timers.target
