# systemd/user/sync-documents.timer
[Unit]
Description=Run sync documents

[Timer]
OnCalendar=*-*-* 00:00:05
WakeSystem=true

[Install]
WantedBy=timers.target
