# systemd/user/sync-documents.service
[Unit]
Description=Sync ~/Documents to Google drive
Wants=network-online.target
After=network-online.target

[Service]
ExecStart=/home/shylock/sh/sync-documents.sh
Type=simple
Restart=on-failure
RestartSec=5

[Install]
WantedBy=default.target
