# systemd/user/sync-tools.timer
[Unit]
Description=Run sync tools

[Timer]
OnCalendar=*-*-* 00:00:10
WakeSystem=true

[Install]
WantedBy=timers.target
