# crontab/jobs
0 2 * * * rclone sync gdrive:/Documents /home/shylock/Documents --transfers 4

